# configs/90-backlight.rules
SUBSYSTEM=="backlight", ACTION=="add", \
  RUN+="/sys/class/backlight/%k/brightness"
